# godaddy-ddns.service
[Unit]
Description=GoDaddy DDNS - update IP address

[Service]
Type=oneshot
ExecStart=/usr/bin/godaddy-ddns --config /etc/godaddy-ddns/config.yaml
